# etc/update-motd.d/41-keymaker-banner
#!/bin/sh
version=$(rpm --query --qf '%{version}-%{release}' keymaker)
cat << 'EOF'

 :::  .   .,::::::.-:.     ::-..        :    :::.      :::  .   .,:::::: :::::::..
 ;;; .;;,.;;;;'''' ';;.   ;;;;';;,.    ;;;   ;;`;;     ;;; .;;,.;;;;'''' ;;;;``;;;;
 [[[[[/'   [[cccc    '[[,[[['  [[[[, ,[[[[, ,[[ '[[,   [[[[[/'   [[cccc   [[[,/[[['
_$$$$,     $$""""      c$$"    $$$$$$$$"$$$c$$$cc$$$c _$$$$,     $$""""   $$$$$$c
"888"88o,  888oo,__  ,8P"`     888 Y88" 888o888   888,"888"88o,  888oo,__ 888b "88bo,
 MMM "MMP" """"YUMMMmM"        MMM  M'  "MMMYMM   ""`  MMM "MMP" """"YUMMMMMMM   "W"

              .::::::..,:::::: :::::::.. :::      .::..,:::::: :::::::..   
             ;;;`    `;;;;'''' ;;;;``;;;;';;,   ,;;;' ;;;;'''' ;;;;``;;;;  
             '[==/[[[[,[[cccc   [[[,/[[[' \[[  .[[/    [[cccc   [[[,/[[['  
               '''    $$$""""   $$$$$$c    Y$c.$$"     $$""""   $$$$$$c    
              88b    dP888oo,__ 888b "88bo, Y88P       888oo,__ 888b "88bo,
               "YMmMY" """"YUMMMMMMM   "W"   MP        """"YUMMMMMMM   "W" 

EOF
echo Version $version
